# scripts/macos/com.saorsa.canvas.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.saorsa.canvas</string>

    <key>ProgramArguments</key>
    <array>
        <string>__INSTALL_DIR__/saorsa-canvas</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>CANVAS_PORT</key>
        <string>9473</string>
        <key>CANVAS_WEB_DIR</key>
        <string>__DATA_DIR__/web</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>__DATA_DIR__/logs/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__DATA_DIR__/logs/stderr.log</string>

    <key>WorkingDirectory</key>
    <string>__DATA_DIR__</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
